# pbr
#!/bin/sh /etc/rc.common

START=99
NFT_FAMILY="inet"
NFT_TABLE="fw4"
MAGIC_IP="10.9.1.2"

start () {
    echo "Starting domain lists update..." >&2

    _process_domain_set \
        "vpn" \
        "https://raw.githubusercontent.com/itdoginfo/allow-domains/refs/heads/main/Categories/geoblock.lst" \
        "https://raw.githubusercontent.com/itdoginfo/allow-domains/refs/heads/main/Services/google_ai.lst"

    _process_domain_set \
        "warp" \
        "https://raw.githubusercontent.com/itdoginfo/allow-domains/refs/heads/main/Categories/block.lst" \
        "https://raw.githubusercontent.com/itdoginfo/allow-domains/refs/heads/main/Categories/porn.lst"  \
        "https://raw.githubusercontent.com/itdoginfo/allow-domains/refs/heads/main/Categories/anime.lst" \
        "https://raw.githubusercontent.com/itdoginfo/allow-domains/refs/heads/main/Categories/news.lst" \
        "https://raw.githubusercontent.com/itdoginfo/allow-domains/refs/heads/main/Services/twitter.lst" \
        "https://raw.githubusercontent.com/itdoginfo/allow-domains/refs/heads/main/Services/meta.lst" \
        "https://raw.githubusercontent.com/labi-le/domains.lst/refs/heads/main/custom-set1.txt"

    nft flush set "$NFT_FAMILY" "$NFT_TABLE" "warp_domains"

    _process_ip_set "warp_domains" "https://core.telegram.org/resources/cidr.txt"
    _process_ip_set "warp_domains" "https://raw.githubusercontent.com/labi-le/domains.lst/refs/heads/main/faf-cf.txt"
    _process_ip_set "warp_domains" "https://raw.githubusercontent.com/routir/unblock/refs/heads/main/services/viber-ip-de.lst"

    echo "Restarting services..." >&2
    /etc/init.d/dnsmasq restart
    /etc/init.d/sing-box restart
    exit 0
}

_process_domain_set() {
    local list_name="$1"
    shift 1
    local urls_list="$@"
    local downloaded_content=""

    echo "Processing $list_name domains..." >&2

    if [ ! -d "/tmp/sing-box/" ]; then
        mkdir -p "/tmp/sing-box/"
    fi

    rm -f "/tmp/dnsmasq.d/00_${list_name}.conf"
    echo '{"version": 1, "rules": []}' > "/tmp/sing-box/${list_name}.json"

    for url in $urls_list; do
        local attempt=1
        local max_attempts=5
        local success=0

        while [ "$attempt" -le "$max_attempts" ]; do
            current_download=$(curl --connect-timeout 5 --max-time 15 -sSLf "$url")
            if [ $? -eq 0 ] && [ -n "$current_download" ]; then
                downloaded_content="${downloaded_content}${current_download}\n"
                success=1
                break
            fi
            echo "  - Error fetching $url. Retrying ($attempt/$max_attempts) in 3s..." >&2
            sleep 3
            attempt=$((attempt + 1))
        done

        if [ "$success" -eq 0 ]; then
             echo "  ! Giving up on $url after $max_attempts attempts." >&2
        fi
    done

    clean_list=$(printf "%b" "$downloaded_content" | awk 'NF>0 && !/^#/{print $1}' | sort -u)

    if [ -n "$clean_list" ]; then
        echo "$clean_list" | awk -v ip="$MAGIC_IP" '{print "address=/" $1 "/" ip}' > "/tmp/dnsmasq.d/00_${list_name}.conf"

        echo "$clean_list" | awk '
            BEGIN { printf "{\"version\": 1, \"rules\": [{\"domain_suffix\": [" }
            { if (NR>1) printf ","; printf "\"%s\"", $1 }
            END { print "]}]}" }
        ' > "/tmp/sing-box/${list_name}.json"

        echo "Successfully generated rules for $list_name." >&2
    fi
}

_process_ip_set() {
    local nft_set_name="$1"
    local ip_url="$2"
    shift 2
    local downloaded_content=""

    echo "Processing IP set: $nft_set_name ($ip_url)" >&2

    local attempt=1
    local max_attempts=5

    while [ "$attempt" -le "$max_attempts" ]; do
        downloaded_content=$(curl --connect-timeout 5 --max-time 15 -sSLf "$ip_url")
        if [ $? -eq 0 ] && [ -n "$downloaded_content" ]; then
            break
        fi
        echo "  - Error fetching IP list. Retrying ($attempt/$max_attempts) in 3s..." >&2
        sleep 3
        attempt=$((attempt + 1))
    done

    if [ -z "$downloaded_content" ]; then
        echo "  ! Giving up on IP list after $max_attempts attempts." >&2
        return 1
    fi

    downloaded_content=$(echo "$downloaded_content" | tr -d '\r' | grep -v ":")
    ip_list_formatted=$(echo "$downloaded_content" | awk 'NF > 0 { printf "%s,", $1 }' | sed 's/,$//')

    if [ -n "$ip_list_formatted" ]; then
        nft add element "$NFT_FAMILY" "$NFT_TABLE" "$nft_set_name" "{ $ip_list_formatted }" 2>/dev/null
    fi
}
