# pbr-yt
#!/bin/sh /etc/rc.common

START=99

start () {
    OUTPUT_FILE="/tmp/dnsmasq.d/youtube.lst"
    > "$OUTPUT_FILE"

    DOMAINS_URLS="
    	https://raw.githubusercontent.com/antonme/ipnames/refs/heads/master/dns-google.txt
        https://raw.githubusercontent.com/antonme/ipnames/refs/heads/master/dns-youtube.txt
        https://raw.githubusercontent.com/GhostRooter0953/discord-voice-ips/refs/heads/master/voice_domains/discord-voice-domains-list
        https://raw.githubusercontent.com/GhostRooter0953/discord-voice-ips/refs/heads/master/main_domains/discord-main-domains-list
    "

    NFT_FAMILY="inet"
    NFT_TABLE="fw4"
    NFT_SET="youtube_domains"

    (
        for DOMAIN_SINGLE_URL in $DOMAINS_URLS; do
            echo "Downloading domain list from $DOMAIN_SINGLE_URL..." >&2
            CURRENT_DOMAIN_DOWNLOAD=$(curl -sSL "$DOMAIN_SINGLE_URL")
            CURL_EXIT_CODE=$?

            if [ "$CURL_EXIT_CODE" -ne 0 ]; then
                echo "Warning: curl failed to download from $DOMAIN_SINGLE_URL (exit code: $CURL_EXIT_CODE). Skipping this URL." >&2
            elif [ -n "$CURRENT_DOMAIN_DOWNLOAD" ]; then
                printf "%s\n" "$CURRENT_DOMAIN_DOWNLOAD"
            else
                echo "Warning: Downloaded content from $DOMAIN_SINGLE_URL is empty." >&2
            fi
        done
    ) | awk -v nft_family="$NFT_FAMILY" -v nft_table="$NFT_TABLE" -v nft_set="$NFT_SET" '
        NF > 0 {
            $1 = $1;
            print "nftset=/" $0 "/4#" nft_family "#" nft_table "#" nft_set
        }
    ' > "$OUTPUT_FILE"

    if [ ! -s "$OUTPUT_FILE" ]; then
        echo "Warning: No domain content was successfully downloaded and processed. Output file ($OUTPUT_FILE) is empty." >&2
    fi

    if dnsmasq --conf-file="$OUTPUT_FILE" --test 2>&1 | grep -q "syntax check OK"; then
        /etc/init.d/dnsmasq restart
    fi

    # update_nft_set_with_cidr
    # if [ $? -ne 0 ]; then
    #     exit 1
    # fi

    exit 0
}


update_nft_set_with_cidr() {
    CIDR_URLS="
        https://raw.githubusercontent.com/touhidurrr/iplist-youtube/refs/heads/main/lists/cidr4.txt
    "

    RAW_CIDR_CONTENT_ACCUMULATOR=""

    for CIDR_SINGLE_URL in $CIDR_URLS; do
        echo "Downloading CIDR list from $CIDR_SINGLE_URL..." >&2
        CURRENT_CIDR_DOWNLOAD=$(curl -sSL "$CIDR_SINGLE_URL")
        CURL_EXIT_CODE=$?

        if [ "$CURL_EXIT_CODE" -ne 0 ]; then
          echo "Error: curl failed to download from $CIDR_SINGLE_URL (exit code: $CURL_EXIT_CODE)." >&2
          return 1
        fi

        if [ -n "$CURRENT_CIDR_DOWNLOAD" ]; then
            RAW_CIDR_CONTENT_ACCUMULATOR="${RAW_CIDR_CONTENT_ACCUMULATOR}${CURRENT_CIDR_DOWNLOAD}"$'\n'
        else
            echo "Warning: Downloaded content from $CIDR_SINGLE_URL is empty." >&2
        fi
    done

    CIDR_CONTENT=$(echo "$RAW_CIDR_CONTENT_ACCUMULATOR" | awk 'NF')

    if [ -z "$CIDR_CONTENT" ]; then
      echo "Error: All downloaded CIDR lists resulted in no valid content." >&2
      return 1
    fi

    total_lines=$(echo "$CIDR_CONTENT" | wc -l | awk '{print $1}')

    if ! [[ "$total_lines" =~ ^[0-9]+$ ]] || [ "$total_lines" -le 0 ]; then
        echo "Error: Could not count lines or downloaded content seems effectively empty (total_lines: $total_lines)." >&2
        return 1
    fi

    ELEMENTS_STRING=$(echo "$CIDR_CONTENT" | awk -v total="$total_lines" '{
      printf "%s", $0;
      if (NR < total) {
        printf ", ";
      }
    }')

    echo "Adding elements to nft set: $NFT_FAMILY/$NFT_TABLE/$NFT_SET..." >&2
    nft add element "$NFT_FAMILY" "$NFT_TABLE" "$NFT_SET" "{ $ELEMENTS_STRING }"
    NFT_EXIT_CODE=$?

    if [ "$NFT_EXIT_CODE" -ne 0 ]; then
      echo "Error: nft add element command failed (exit code: $NFT_EXIT_CODE)." >&2
      return 1
    fi
    return 0
}
